# jehon-base-minimal/usr/share/jehon-base-minimal/etc/cron.daily/jh-backup-computer
#!/usr/bin/env bash

jh-backup-computer
